# debian/postinst
#!/bin/bash
echo "Enabling and starting joycond service..."
systemctl enable --now joycond